# logstash-forwarder/logstash-forwarder.upstart.conf
start on runlevel [2345]
stop on runlevel [06]

console log
umask 0022

script
  /usr/bin/logstash-forwarder -config /etc/logstash-forwarder -spool-size 100 -log-to-syslog
end script
